# deployment/systemd/skulk.service
[Unit]
Description=Skulk distributed inference node
Documentation=https://github.com/Foxlight-Foundation/Skulk
# Network-online ensures DNS / outbound connectivity is available before
# Skulk attempts libp2p bootstrap. local-fs ensures $HOME/.skulk and the
# model store paths are mounted (matters on systems with delayed mounts).
After=network-online.target local-fs.target
Wants=network-online.target
# StartLimitBurst/StartLimitIntervalSec are unit-level directives — they
# must live in [Unit], not [Service], or systemd silently ignores them.
# 5 crashes within 5 minutes trips the limit and stops restart attempts.
StartLimitBurst=5
StartLimitIntervalSec=300

[Service]
Type=simple
WorkingDirectory=__SKULK_REPO__

# Operator env file. EnvironmentFile loads KEY=VALUE pairs into the
# unit's environment before ExecStart, matching the macOS LaunchAgent's
# wrapper-sourced behavior. Leading `-` makes the file optional so a
# fresh install (before the file is copied) doesn't fail to start.
EnvironmentFile=-%h/.skulk/skulk.env

# Boot-time prep + skulk launch is delegated to the wrapper so the
# failure policy (best-effort git pull / uv sync, hard requirement on
# a built dashboard) is identical to the macOS install. The wrapper
# also sources ~/.skulk/skulk.env again so values inside the env file
# (like SKULK_AUTO_UPDATE) reach uv / npm subprocesses.
ExecStart=__SKULK_REPO__/deployment/install/skulk-startup.sh

# SIGTERM is what `Skulk._signal` already handles via task-group cancel.
# `mixed` sends SIGTERM to the main process and SIGKILL to remaining
# children only if they outlive TimeoutStopSec — important because Skulk
# spawns runner subprocesses that hold MLX state.
KillSignal=SIGTERM
KillMode=mixed
TimeoutStopSec=30

# Crash recovery: only restart on abnormal exit. Planned shutdowns
# (`systemctl --user stop skulk` → SIGTERM → exit 0) stay stopped.
# `RestartSec` avoids a tight loop when the cause is a config error;
# the [Unit] StartLimitBurst/StartLimitIntervalSec bail after repeated
# failure so a broken config doesn't burn the node.
Restart=on-failure
RestartSec=10

# MLX engines and the libp2p stack each open many file descriptors and
# memory-map model shards. The defaults (1024 / system rlimit) are too
# tight for a node hosting multi-GB models.
LimitNOFILE=65536
LimitMEMLOCK=infinity

# Logging: stdout/stderr flow to the journal. The Vector sidecar
# (deployment/logging/vector.yaml) is responsible for shipping the
# JSON-on-stdout stream when centralized logging is enabled.
StandardOutput=journal
StandardError=journal
SyslogIdentifier=skulk

[Install]
WantedBy=default.target
