# deployment/launchd/foundation.foxlight.skulk.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>foundation.foxlight.skulk</string>

    <!-- Boot-time prep + skulk launch is delegated to a wrapper script
         so the failure policy (best-effort git pull / uv sync, hard
         requirement on a built dashboard) can be edited without
         re-installing the agent. The wrapper sources ~/.skulk/skulk.env
         for operator-configurable env vars.

         The wrapper is exec'd directly (not via `bash -lc`) so a repo
         path containing spaces or shell metacharacters does not break
         tokenization. PATH and other env vars are supplied via the
         EnvironmentVariables block below. -->
    <key>ProgramArguments</key>
    <array>
        <string>__SKULK_REPO__/deployment/install/skulk-startup.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__SKULK_REPO__</string>

    <!-- Run at user login (LaunchAgent) and keep alive on crash, but
         not on clean exit. SuccessfulExit=false means launchd treats
         exit 0 as "stay stopped" and any non-zero as "relaunch".
         Crashed=true also relaunches on signal-induced termination. -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <!-- Throttle relaunch to avoid a tight loop when the cause is a
         config error. Matches the systemd unit's RestartSec=10. -->
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <!-- Interactive ProcessType is the default for LaunchAgents and
         ensures Metal / window-server access works the same way it
         does for `uv run skulk` from Terminal.app. -->
    <key>ProcessType</key>
    <string>Interactive</string>

    <!-- Captured stdout/stderr from skulk + the wrapper's exec. The
         wrapper also writes a separate prep log at
         ~/.skulk/logs/skulk.prep.log so operators can audit boot-time
         updates independently of skulk's own log stream. -->
    <key>StandardOutPath</key>
    <string>__SKULK_LOG_DIR__/skulk.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>__SKULK_LOG_DIR__/skulk.stderr.log</string>

    <!-- Inherit the user's shell PATH so `uv`, `npm`, `git`, and any
         tool-version shims resolve. All other env vars come from
         ~/.skulk/skulk.env (sourced by the wrapper) so they can be
         edited without re-installing the agent. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__USER_PATH__</string>
    </dict>
</dict>
</plist>
